# publicip.cron
*/1 * * * * /path/to/executable/script && ./script.sh
